# mysql/mysql.conf
[program:mysql]
command=/etc/supervisor/conf.d/mysql.sh
redirect_stderr=true
